# roles/xray/templates/raven-subscribe/xray-subscription.service.j2
[Unit]
Description=Xray Subscription Server (Raven-subscribe)
Documentation=https://github.com/AlchemyLink/Raven-subscribe
After=network.target {{ xray_service_name }}.service

[Service]
Type=simple
User={{ xray_user }}
Group={{ xray_group }}
ExecStart={{ raven_subscribe_install_dir }}/xray-subscription -config {{ raven_subscribe_config_dir }}/config.json
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ raven_subscribe_db_dir }} {{ xray_config_dir }}config.d
ReadOnlyPaths={{ raven_subscribe_config_dir }}
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
